# service-registration-and-discovery-service/app.service.j2
[Unit]
Description={{ app_name }}
After=syslog.target

[Service]
User={{ username }}
ExecStart="{{ jar_target }}"
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target